# packages/bsp/common/etc/cron.d/atrios-updates
@reboot root /usr/lib/atrios/atrios-apt-updates
@daily root /usr/lib/atrios/atrios-apt-updates
