# packages/bsp/atrios-lowmem/lowmem-mkswap.service
[Unit]
Description=AtriOS lowmem: Create a swapfile
DefaultDependencies=no
Wants=atrios-resize-filesystem.service
After=atrios-resize-filesystem.service
Before=basic.target

[Service]
Type=oneshot
ExecStart=/usr/bin/lowmem-mkswap.sh
RemainAfterExit=yes

[Install]
WantedBy=basic.target
