# packages/bsp/common/lib/systemd/system/atrios-firstrun.service
# AtriOS firstrun service
# Regenerates the SSH keys and performs additional first run tweaks
# This service will run in parallel with other services

[Unit]
Description=AtriOS first run tasks
Before=getty.target system-getty.slice
After=ssh.service

[Service]
Type=simple
RemainAfterExit=yes
EnvironmentFile=/etc/default/atrios-firstrun
ExecStart=/usr/lib/atrios/atrios-firstrun start
TimeoutStartSec=2min

[Install]
WantedBy=multi-user.target

